# deauwait.service
[Unit]
Description=DeauWait — YouTube Live Watcher (Nuxt)
After=network.target

[Service]
Type=simple
User=noble
WorkingDirectory=/home/noble/deauwait
# Node 20.6+ supports --env-file natively, no dotenv package needed
ExecStart=/usr/bin/node --env-file=.env .output/server/index.mjs
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
